# srcs/grafana/srcs/supervisord.conf
[supervisord]
nodaemon=true
user=root

[program:grafana-server]
command=sh -c "kill $(pidof grafana-server) || /usr/sbin/grafana-server --config=/etc/grafana.ini --homepath /usr/share/grafana/ && kill $(pidof supervisord)"
autostart=true
autorestart=true

[program:telegraf]
command=sh -c "kill $(pidof telegraf) || telegraf && kill $(pidof supervisord)"
autostart=true
autorestart=true